# interpres.service.example
# Copy this file to interpres.service and replace YOUR_USERNAME with your
# actual username, then:
#   systemctl --user enable interpres.service
#   systemctl --user start interpres.service

[Unit]
Description=Interpres Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/YOUR_USERNAME/projects/interpres
ExecStart=/home/YOUR_USERNAME/projects/interpres/venv/bin/python3 /home/YOUR_USERNAME/projects/interpres/bot.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/YOUR_USERNAME/projects/interpres/.env
Environment=PATH=/home/YOUR_USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=DISPLAY=:1

[Install]
WantedBy=default.target
